# jenkins/root/etc/supervisor/conf.d/jenkins.conf
[program:jenkins]
command=/usr/local/bin/start-jenkins.sh
stdout_events_enabled=true
stderr_events_enabled=true
